# packaging/systemd/ai-r-mcp.service
[Unit]
Description=ai-r MCP shared server (streamable-http, socket-activated)
Documentation=https://github.com/dmakarov/ai-r
Requires=ai-r-mcp.socket
After=ai-r-mcp.socket

[Service]
Type=simple
# http transport; the server picks up the socket systemd passes via
# LISTEN_FDS/LISTEN_PID (ai_r.serve.systemd_listen_sockets).
Environment=AI_R_MCP_TRANSPORT=http
# Self-exit after 15 min idle (never mid-request). systemd keeps the socket
# and re-activates the service on the next connection -> zero resident
# process when idle, respawn on demand.
Environment=AI_R_MCP_IDLE_SEC=900
ExecStart=%h/.local/bin/ai-r-mcp
# Idle self-exit is normal, not a failure; let the socket (not systemd) drive
# respawn, so do not auto-restart on clean exit.
Restart=no

[Install]
# Reached through the socket; no direct WantedBy needed.
